# agent/buildos-agent.service
[Unit]
Description=BuildOS Infra Remote Agent
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=buildos-agent
Group=docker
WorkingDirectory=/opt/buildos-agent
EnvironmentFile=/etc/buildos-agent/config.env
ExecStart=/opt/buildos-agent/venv/bin/python /opt/buildos-agent/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
